# .docker/configs/etc/supervisor/conf.d/php.conf
[group:php-server]
programs=php-server
priority=20

[program:php-server]
command=php /app/server.php
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/app/php.log
